# meta-opencentauri/recipes-core/busybox/files/crontab
*/15 * * * * /usr/sbin/ntpd -nq -p pool.ntp.org
